# dvr-rtsp.service
[Unit]
Description=DVR RTSP Server (mediamtx + HiEasy DVR bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dvr
Group=dvr
WorkingDirectory=/opt/dvr
ExecStart=/opt/dvr/mediamtx /opt/dvr/mediamtx.yml
Restart=always
RestartSec=5

# DVR connection settings — edit /opt/dvr/dvr.env to match your DVR
EnvironmentFile=/opt/dvr/dvr.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dvr-rtsp

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/dvr

[Install]
WantedBy=multi-user.target
